# systemd/rldc-api.service
[Unit]
Description=RLdC Trading Bot FastAPI Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/rldc-trading-bot
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONPATH=/var/www/rldc-trading-bot"
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=append:/var/log/rldc-api.log
StandardError=append:/var/log/rldc-api.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
